# systemd/bpf-transfer.service
[Unit]
Description=Start bpf-transfer

[Service]
ExecStart=/root/server.py

[Install]
WantedBy=multi-user.target
